# docker/worker/supervisord.conf
[supervisord]
nodaemon=true
user=root
logfile=/var/www/storage/logs/supervisord.log
logfile_maxbytes=10MB
pidfile=/var/run/supervisord.pid

; ─── Reader: lê o CSV e distribui chunks ────────────────────────────────────
[program:import-reader]
process_name=import_reader_%(process_num)02d
numprocs=%(ENV_QUEUE_READERS)s
command=php /var/www/artisan importacao:reader
user=www-data
autostart=true
autorestart=true
startsecs=2
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/worker.log
stdout_logfile_maxbytes=20MB

; ─── Workers: validam chunks, gravam spool CSV e JSONL de erros ─────────────
[program:import-worker]
process_name=import_worker_%(process_num)02d
numprocs=%(ENV_QUEUE_WORKERS)s
command=php /var/www/artisan importacao:worker
user=www-data
autostart=true
autorestart=true
startsecs=2
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/worker.log
stdout_logfile_maxbytes=20MB

; ─── Spool: carrega CSV no Postgres via COPY ─────────────────────────────────
[program:import-spool]
process_name=import_spool_%(process_num)02d
numprocs=%(ENV_QUEUE_SPOOL)s
command=php /var/www/artisan importacao:spool
user=www-data
autostart=true
autorestart=true
startsecs=2
redirect_stderr=true
stdout_logfile=/var/www/storage/logs/worker.log
stdout_logfile_maxbytes=20MB
